# examples/ydbd/loglugger-client.service
[Unit]
Description=Loglugger client for YDBD journald logs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/loglugger-client
ExecStart=/opt/loglugger/bin/client \
  -server ${LOGLUGGER_SERVER_URL} \
  -client-id ${LOGLUGGER_CLIENT_ID} \
  -service-mask ${LOGLUGGER_SERVICE_MASK} \
  -message-regex '^(?P<P_DTTM>[^ ]*) :(?P<P_SERVICE>[^ ]*) (?P<P_LEVEL>[^ ]*): (?P<P_MESSAGE>.*).*$' \
  -message-regex-no-match send_raw \
  -batch-size ${LOGLUGGER_BATCH_SIZE} \
  -batch-timeout ${LOGLUGGER_BATCH_TIMEOUT} \
  -http-timeout ${LOGLUGGER_HTTP_TIMEOUT} \
  -retry-delay ${LOGLUGGER_RETRY_DELAY} \
  -tls-ca-file ${LOGLUGGER_TLS_CA_FILE} \
  -tls-cert-file ${LOGLUGGER_TLS_CERT_FILE} \
  -tls-key-file ${LOGLUGGER_TLS_KEY_FILE}
Restart=always
RestartSec=2s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
